# misc/linux/freedom-routes.timer
[Unit]
Description=Freedom Routes

[Timer]
OnCalendar=weekly

[Install]
WantedBy=default.target
