# assets/linux/systemd/freedom-routes.timer
[Unit]
Description=Freedom Routes

[Timer]
OnCalendar=daily

[Install]
WantedBy=default.target
